# pimacrkeys.service
[Unit]
Description=Pi Macro Screen - USB HID Keyboard Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/piMacroScreen
ExecStart=/usr/bin/python3 /home/pi/piMacroScreen/macroKeys.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
