# pimacrkeys.service
[Unit]
Description=Pi Macro Screen - USB HID Keyboard Service
After=network.target

[Service]
Type=simple
User=root
Environment="APP_DIR=/home/lheringer/piMacroScreen"
WorkingDirectory=/home/lheringer/piMacroScreen
ExecStart=/home/lheringer/piMacroScreen/start_services.sh start
ExecStop=/home/lheringer/piMacroScreen/start_services.sh stop
ExecReload=/home/lheringer/piMacroScreen/start_services.sh reload
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
